# assets/root/etc/init.d/tor_settings
#!/sbin/openrc-run

description="Sets up tor settings from envoriment"

start() {
    ebegin "Setting up Tor settings"
    bash /tor.sh
    eend $?
}
